# services/cleanup.service
[Unit]
Description=Bird Detection Storage Cleanup
After=network.target

[Service]
Type=oneshot
User=jaysettle
Group=jaysettle
WorkingDirectory=/home/jaysettle/JayModel/MotionRevamp
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -c "from src.cleanup_manager import run_cleanup_service; run_cleanup_service()"

[Install]
WantedBy=multi-user.target